# phase_ii/systemd/timers/monthly_summary.timer
[Unit]
Description=Run Monthly Snowflake Summary Log

[Timer]
OnCalendar=*-*-01 00:25:00
Persistent=true

[Install]
WantedBy=timers.target
