# phase_ii/systemd/timers/weekly_summary.timer
[Unit]
Description=Run Weekly Snowflake Summary Log

[Timer]
OnCalendar=Sun *-*-* 00:15:00
Persistent=true

[Install]
WantedBy=timers.target
